# deploy/airspace-monitor.service
[Unit]
Description=Raspberry Pi Airspace Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-airspace-monitor
EnvironmentFile=-/etc/airspace-monitor.env
ExecStart=/usr/bin/python3 -m airspace_monitor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
